# rpms/BUILD/lockss-daemon-root/etc/init.d/lockss
#!/bin/sh

# Start and stop LOCKSS daemon
#
# chkconfig: 345 99 01
# description: LOCKSS daemon

. /etc/init.d/functions
. /etc/profile
. /etc/lockss/functions

do_command () {
if [ -f ${LOCKSS_HOME}/config.$2.dat ]
then
	${CFG_FILE}=${LOCKSS_HOME}/config.$2.dat
fi
case "$1" in
    start)
	if [ ! -r ${CFG_FILE} ]; then
	    echo "Can't start ${2}: ${CFG_FILE} missing or not readable"
	    exit 1
	fi
	if [ ! -d ${LOG_DIR} ]; then
	    mkdir -p ${LOG_DIR}
	    chown $2 ${LOG_DIR}
	    chmod 755 ${LOG_DIR}
	fi
	if [ ! -f ${PID_FILE} ]; then
	    touch ${PID_FILE}
	fi
	chown $2 ${PID_FILE}
	chmod 644 ${PID_FILE}
	echo -n "Starting ${2}: "
	echo "Starting LOCKSS for user $2 at `date`" >> ${LOG_FILE}
	chown $2 ${LOG_FILE}
	daemon /etc/lockss/startdaemon $2
	echo
	;;
    stop)
	echo -n "Stopping ${2}: "
	rm -f "${KEEP_GOING}"
	killproc $2
	echo
	;;
    *)
	echo "Usage: $0 {start|stop}"
	exit 1
	;;
esac
}

set_lockss_user
for A in ${LOCKSS_USER}
do
        set_variables ${A}
        do_command $1 $A
	unset_variables
done


exit 0
