# rpms/BUILD/lockss-daemon-root/etc/init.d/lockss
#!/bin/sh

# Start and stop LOCKSS daemon
#
# chkconfig: 345 99 01
# description: LOCKSS daemon

. /etc/init.d/functions
. /etc/profile
. /etc/lockss/functions

do_command () {
if [ -f /etc/lockss/config.$2.dat ]
then
        CFG_FILE=/etc/lockss/config.$2.dat
fi
case "$1" in
    start)
	start_lockss ${2}
        ;;
    stop)
        stop_lockss ${2}
        ;;
    restart)
	stop_lockss ${2}
	start_lockss ${2}
	;;
    killall)
        echo "killall is no longer supported, sorry.
        exit 1
#         killall_lockss ${2}
        ;;
    *)
        echo "Usage: $0 {start|stop|restart}"
        exit 1
        ;;
esac
}

stop_lockss () {
    echo -n "Stopping ${1}: "
    echo "Stopping LOCKSS for user ${1} at `date`" >> ${LOG_FILE}
    rm -f "${KEEP_GOING}"
    killproc ${1}

    count=0
    while [ -e ${LOCKFILE} -a $count -lt 15 ] ; do
        sleep 1
        count=`expr $count + 1`
    done

    rm -f ${PID_FILE}
    # tk - should kill ssl proc, del ssl  passwd file and pid file
    echo
}

start_lockss () {
    if [ ! -r ${CFG_FILE} ]; then
        echo "Can't start ${1}: ${CFG_FILE} missing or not readable"
        exit 1
    fi
    if [ ! -d ${LOG_DIR} ]; then
        mkdir -p ${LOG_DIR}
        chown ${1} ${LOG_DIR}
        chmod 755 ${LOG_DIR}
    fi
    if [ ! -f ${PID_FILE} ]; then
        touch ${PID_FILE}
    fi
    chown ${1} ${PID_FILE}
    chmod 644 ${PID_FILE}
    echo -n "Starting ${1}: "
    echo "Starting LOCKSS for user ${1} at `date`" >> ${LOG_FILE}
    chown ${1} ${LOG_FILE}
    daemon /etc/lockss/startdaemon ${1}
    echo
}

killall_lockss () {
    echo -n "Killing all ${1} processes: "
    rm -f "${KEEP_GOING}" "${LOCKFILE}" "${PID_FILE}"
    killall -u $1 -s 9

    # tk - should kill ssl proc, del ssl  passwd file and pid file
    echo
}

set_lockss_user
for A in ${LOCKSS_USER}
do
        set_variables ${A}
        do_command ${1} $A
        unset_variables
done


exit 0
